# eu.openanalytics.architect.r.server.win32.win32.x86_64/build_tools/rtools40/etc/profile.d/gitbin.sh
# Set path to Git for Windows
[ -d "/c/progra~1/git/cmd" ] && PATH="$PATH:/c/progra~1/git/cmd"
export PATH
